# deploy/crackthedeck-backend.service
[Unit]
Description=CrackTheDeck API (FastAPI)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/crackthedeck/crackthedeck-backend/crackthedeck-backend
# Backend .env first, then project .env in CTD/crackthedeck (overrides) — /var/www/crackthedeck/.env
EnvironmentFile=/var/www/crackthedeck/crackthedeck-backend/crackthedeck-backend/.env
EnvironmentFile=-/var/www/crackthedeck/.env
ExecStart=/var/www/crackthedeck/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
